# usb-gadget.service
[Unit]
Description=USB Gadget Setup
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb-keyboard-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
